# supervisor-etcd/supervisor-etcd.conf
[program:etcd]
command=/usr/local/bin/etcd -addr 127.0.0.1 -peer-addr 127.0.0.1
autostart=true
autorestart=true
startsecs=1
